# controllers/20-ds4.rules
SUBSYSTEM=="input", ACTION=="add", ATTRS{name}=="Wireless Controller", RUN+="/bin/bash /home/pi/RetroPie/myscripts/ds4LedColor.sh"
